# projects/musicbox/rootfs/etc/rc.local
#!/bin/sh

if [ -x /root/rc.local ]; then
	/root/rc.local
fi
