# projects/musicbox/rootfs/var/spool/cron/crontabs/root
# MIN HOUR DAY MONTH DAYOFWEEK CMD

0 7 * * 1-5 /root/wakeup.sh
17 4 * * * /root/ntp.sh
