# projects/dvb/rootfs/var/spool/cron/crontabs/root
# MIN HOUR DAY MONTH DAYOFWEEK CMD

17 4 * * * /root/ntp.sh
