# projects/dvb/rootfs/etc/udev/rules.d/81-ttyusb.rules
SUBSYSTEM=="tty", KERNEL=="ttyGS?", RUN+="/bin/runc %E{ACTION} getty %k 0 115200"
